# files/Chapter-10/python/iloop.service
[Unit]
Description=

[Service]
Type=
PIDfile=/run/iloop.pid
ExecStart=
PrivateTmp=

[Install]
WantedBy=
